# container/supervisord.conf
[supervisord]
nodaemon=true

[program:novnc]
command=/home/user/noVNC/utils/launch.sh --vnc localhost:5900
startretries=0
